# examples/launchd/com.user.slurm-status.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
    LaunchAgent to auto-start Slurm status monitoring script on macOS login

    Installation:
    1. Copy this file to ~/Library/LaunchAgents/
    2. Rename it (e.g., com.myusername.slurm-status.plist)
    3. Update USER, the environment variables, and the cluster placeholders below
    4. Update script path if different
    5. Load with:
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.myusername.slurm-status.plist
       launchctl enable gui/$(id -u)/com.myusername.slurm-status
       launchctl kickstart -k gui/$(id -u)/com.myusername.slurm-status
    -->

    <key>Label</key>
    <string>com.user.slurm-status</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- Update this path to where you cloned the repo -->
        <string>/Users/USER/vscode-slurm-status-bar/scripts/slurm_monitor.sh</string>
        <!-- Update these to the logical remotes you want to display -->
        <string>fir</string>
        <string>ror</string>
        <string>nibi</string>
        <string>tril</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- LaunchAgent PATH is minimal on macOS; include common CLI locations -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>

        <!-- Required when your local username differs from your cluster username -->
        <key>SLURM_STATUS_BAR_REMOTE_USER</key>
        <string>YOUR_CLUSTER_USERNAME</string>

        <!-- Prefer automation-node aliases for unattended use -->
        <key>SLURM_STATUS_BAR_DISABLE_PASSCODE_FALLBACK</key>
        <string>1</string>
        <key>SLURM_STATUS_BAR_DISABLE_FAIRSHARE</key>
        <string>1</string>

        <key>SLURM_STATUS_BAR_SQUEUE_ALIAS_FIR</key>
        <string>robot-fir-slurm</string>
        <key>SLURM_STATUS_BAR_SQUEUE_ALIAS_ROR</key>
        <string>robot-ror-slurm</string>
        <key>SLURM_STATUS_BAR_SQUEUE_ALIAS_NIBI</key>
        <string>robot-nibi-slurm</string>
        <key>SLURM_STATUS_BAR_SQUEUE_ALIAS_TRIL</key>
        <string>robot-tril-slurm</string>

        <!-- Fairshare aliases can be added later if your automation-node setup gets an approved sshare wrapper -->
        <!--
        <key>SLURM_STATUS_BAR_FAIRSHARE_ALIAS_FIR</key>
        <string>robot-fir-fairshare</string>
        <key>SLURM_STATUS_BAR_FAIRSHARE_ALIAS_ROR</key>
        <string>robot-ror-fairshare</string>
        <key>SLURM_STATUS_BAR_FAIRSHARE_ALIAS_NIBI</key>
        <string>robot-nibi-fairshare</string>
        -->
    </dict>

    <!-- Start automatically on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep the script running (restart if it crashes) -->
    <key>KeepAlive</key>
    <true/>

    <!-- Optional: Log output for debugging -->
    <key>StandardOutPath</key>
    <string>/tmp/slurm-status.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/slurm-status.err</string>
</dict>
</plist>
